# deploy/systemd/kipclip.service
[Unit]
Description=kipclip Deno app (Fresh)
After=network.target

[Service]
Type=simple
User=kipclip
Group=kipclip
WorkingDirectory=/var/lib/kipclip/current
EnvironmentFile=/etc/kipclip/env
ExecStart=/opt/deno/bin/deno serve -A --host 127.0.0.1 --port 8000 _fresh/server.js
Restart=always
RestartSec=2s
# Deno waits for long-lived WebSocket connections during graceful shutdown.
# Bound release/reboot downtime before systemd force-closes them.
TimeoutStopSec=10s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/kipclip /var/log/kipclip

[Install]
WantedBy=multi-user.target
