# deploy/systemd/restic-backup.timer
[Unit]
Description=Run kipclip restic backup nightly at 04:00 UTC

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
