# deploy/systemd/tap.service
[Unit]
Description=TAP firehose subscriber (relay -> kipclip webhook)
After=network.target

[Service]
Type=simple
User=tap
Group=tap
WorkingDirectory=/var/lib/tap
EnvironmentFile=/etc/tap/env
ExecStart=/opt/tap/tap run
Restart=always
RestartSec=5s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/tap /var/log/tap

[Install]
WantedBy=multi-user.target
